# src/agent/service/wazuh-agent.service
[Unit]
Description=Wazuh agent
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=simple
PIDFile=/var/run/wazuh-agent.lock
ExecStart=/usr/bin/env WAZUH_HOME/wazuh-agent
TimeoutStopSec=30s    # Wait for 30 seconds before killing the service

KillSignal=SIGTERM

KillMode=mixed

SendSIGKILL=no

RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
